# scripts/systemd/ssh-tunnel.service.d/override.conf
[Service]
Environment="IDENTITY_FILE=/home/pi/.ssh/ohio2.pem"
Environment="DESTINATION=ohio@54.237.238.250"
Environment="BIND_PORT=5003"
